# sdrservice.conf
description "My rtl_tcp service"
author "your@email-address.com"

start on runlevel [2345]
stop on runlevel [!2345]

start on started networking

respawn

script
  sleep 30
  exec /home/pi/rtl_tcp.sh
end script
